# deploy/ceo-dashboard.service
[Unit]
Description=CEO Executive Dashboard (Next.js standalone)
After=network-online.target
Wants=network-online.target

[Service]
Type=simple
User=ceo
Group=ceo
WorkingDirectory=/opt/ceo-dashboard
EnvironmentFile=/opt/ceo-dashboard/.env.production
# Next standalone server; bind to localhost — Caddy terminates TLS in front.
Environment=PORT=3000
Environment=HOSTNAME=127.0.0.1
Environment=NODE_ENV=production
ExecStart=/usr/bin/node server.js
Restart=always
RestartSec=5
MemoryMax=1G
StandardOutput=journal
StandardError=journal
SyslogIdentifier=ceo-dashboard

[Install]
WantedBy=multi-user.target
