# deploy/lexi-hermes.service
[Unit]
Description=Lexi headless worker (orchestrator + Composio webhook ingress)
After=network-online.target
Wants=network-online.target

[Service]
Type=simple
User=lexi
Group=lexi
WorkingDirectory=/opt/lexi
Environment=LEXI_ENV=production
Environment=LEXI_ENV_FILE=/opt/lexi/.env.production
ExecStart=/opt/lexi/.venv/bin/python -m app.worker --webhook
Restart=always
RestartSec=5
# Resource guardrails for the KVM2 VPS.
MemoryMax=1G
# Log to the journal; file logs rotate via RotatingFileHandler in-process.
StandardOutput=journal
StandardError=journal
SyslogIdentifier=lexi-hermes

[Install]
WantedBy=multi-user.target
