# deploy/lexi-backup.timer
[Unit]
Description=Run Lexi SQLite backup hourly

[Timer]
OnCalendar=hourly
Persistent=true

[Install]
WantedBy=timers.target
